# home/dot_config/powershell/profile.d/00-aliases.ps1.tmpl
# Configure aliases
Set-Alias -Name c -Value clear
Set-Alias -Name cc -Value claude
Set-Alias -Name cm -Value chezmoi
Set-Alias -Name cmm -value chezmoi_modify_manager
Set-Alias -Name g -Value git
Set-Alias -Name htop -Value ntop
{{- if eq .chezmoi.os "windows" }}
# `ls` is a default alias only on Windows; drop it so the eza `ls` function
# (10-functions) wins. On Unix `ls` is the native binary, nothing to remove.
Remove-Alias -Name ls
{{- end }}
Set-Alias -Name pn -Value pnpm
{{- if and (eq .hosttype "ewn") (eq .chezmoi.os "windows") }}
Set-Alias -Name tg -Value TortoiseGitProc
{{- end }}
Set-Alias -Name top -Value ntop
Set-Alias -Name v -Value vim
